# deploy/systemd/flydb.service
[Unit]
Description=FlyDB - High-Performance In-Memory Database
Documentation=https://github.com/flydb/flydb
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=flydb
Group=flydb

# Working directory
WorkingDirectory=/var/lib/flydb

# Binary and configuration
ExecStart=/usr/local/bin/flydb --config /etc/flydb/flydb.toml

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/flydb
ReadOnlyPaths=/etc/flydb

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flydb

# Environment
Environment="FLYDB_LOG_LEVEL=info"

[Install]
WantedBy=multi-user.target

